# doh-hole.service
[Unit]
Description=DoH Hole

[Service]
ExecStart=/usr/bin/doh-hole

[Install]
WantedBy=multi-user.target